# etc/systemd/system/shell-php.service
[Unit]
Description=Shell PHP

[Service]
ExecStart=echo "hello world"

[Install]
WantedBy=multi-user.target